# babysit-arc/arc-babysitter.service
# ARC Babysitter — systemd USER unit (reboot-resilient supervisor for the babysit-arc skill).
#
# Auto-starts the poll-loop supervisor on boot and rebuilds campaign state from STATUS.md, so a power
# outage / reboot self-recovers an unattended ARC campaign. The supervisor's own lockfile still
# guarantees a single orchestrator.
#
# Install (user scope — no root):
#   mkdir -p ~/.config/systemd/user
#   cp ~/.claude/skills/babysit-arc/arc-babysitter.service ~/.config/systemd/user/
#   loginctl enable-linger "$USER"           # let the user service run without an active login
#   systemctl --user daemon-reload
#   systemctl --user enable --now arc-babysitter
#   journalctl --user -u arc-babysitter -f   # or: tail -f ~/Projects/arc_babysitter.log
#
# Stop / disable:
#   systemctl --user stop arc-babysitter
#   systemctl --user disable arc-babysitter
#
# NB: the supervisor calls `claude -p`, which needs `claude` (and conda for the agent's `conda run`)
# on PATH — so we run it through a LOGIN shell (`bash -lc`) to source ~/.profile / ~/.bashrc / conda
# init. If your PATH/conda setup lives elsewhere, adjust the ExecStart line.

[Unit]
Description=ARC Babysitter supervisor (babysit-arc skill)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash -lc 'exec "%h/.claude/skills/babysit-arc/arc_babysitter.sh"'
Restart=on-failure
RestartSec=60
# exit 0 = pool DONE (clean stop); exit 3 = another supervisor holds the lock (don't fight it).
RestartPreventExitStatus=3

[Install]
WantedBy=default.target
